# TP-5/Yt/yt.timer
[Unit]
Description=Timer pour le service YT

[Timer]
OnCalendar=minutely
Unit=yt.service

[Install]
WantedBy=timers.target